# supervisord-fly.conf
# supervisord configuration for Fly.io single-VM deployment.
# Starts all backend Deno services + Traefik reverse proxy + static frontend + Redpanda.
# All services communicate on localhost with fixed ports.

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
user=root

[program:traefik]
command=/usr/local/bin/traefik --configFile=/app/traefik-fly.yml
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redpanda]
command=bash -c "mkdir -p /var/lib/redpanda/data /etc/redpanda && /usr/local/bin/rpk redpanda start --mode=dev-container --smp=1 --memory=200M --default-log-level=warn --kafka-addr=0.0.0.0:9092 --advertise-kafka-addr=localhost:9092"
autostart=true
autorestart=true
priority=2
startsecs=8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redpanda-console]
command=/usr/local/bin/redpanda-console
environment=KAFKA_BROKERS="localhost:9092",SERVER_LISTENPORT="8081"
autostart=true
autorestart=true
priority=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:market-sim]
command=deno run --allow-all /app/backend/src/market-sim/market-sim.ts
environment=MARKET_SIM_PORT="5000",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:observability]
command=deno run --allow-all /app/backend/src/observability/observability-server.ts
environment=OBSERVABILITY_PORT="5007",OBS_DB_PATH="/app/backend/data/observability.db",OBS_RETENTION_MS="86400000",REDPANDA_BROKERS="localhost:9092",MARKET_SIM_PORT="5000",EMS_PORT="5001",OMS_PORT="5002",ALGO_TRADER_PORT="5003",TWAP_ALGO_PORT="5004",POV_ALGO_PORT="5005",VWAP_ALGO_PORT="5006",USER_SERVICE_PORT="5008",JOURNAL_PORT="5009"
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ems]
command=deno run --allow-all /app/backend/src/ems/ems-server.ts
environment=EMS_PORT="5001",MARKET_SIM_HOST="localhost",MARKET_SIM_PORT="5000",EMS_PARTICIPATION_CAP="0.20",EMS_IMPACT_PER_1000_BPS="1.0",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fix-exchange]
command=deno run --allow-all /app/backend/src/fix/fix-exchange.ts
environment=FIX_EXCHANGE_PORT="9880",MARKET_SIM_HOST="localhost",MARKET_SIM_PORT="5000",EMS_PARTICIPATION_CAP="0.20",EMS_IMPACT_PER_1000_BPS="1.0"
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fix-gateway]
command=deno run --allow-all /app/backend/src/fix/fix-gateway.ts
environment=FIX_GATEWAY_PORT="9881",FIX_EXCHANGE_HOST="localhost"
directory=/app
autostart=true
autorestart=true
priority=25
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:oms]
command=deno run --allow-all /app/backend/src/oms/oms-server.ts
environment=OMS_PORT="5002",ALGO_HOST="localhost",ALGO_TRADER_PORT="5003",TWAP_ALGO_PORT="5004",POV_ALGO_PORT="5005",VWAP_ALGO_PORT="5006",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:limit-algo]
command=deno run --allow-all /app/backend/src/algo/limit-strategy.ts
environment=ALGO_TRADER_PORT="5003",EMS_HOST="localhost",EMS_PORT="5001",MARKET_SIM_HOST="localhost",MARKET_SIM_PORT="5000",OBSERVABILITY_HOST="localhost",OBSERVABILITY_PORT="5007",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:twap-algo]
command=deno run --allow-all /app/backend/src/algo/twap-strategy.ts
environment=TWAP_ALGO_PORT="5004",EMS_HOST="localhost",EMS_PORT="5001",MARKET_SIM_HOST="localhost",MARKET_SIM_PORT="5000",TWAP_INTERVAL_MS="5000",OBSERVABILITY_HOST="localhost",OBSERVABILITY_PORT="5007",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pov-algo]
command=deno run --allow-all /app/backend/src/algo/pov-strategy.ts
environment=POV_ALGO_PORT="5005",EMS_HOST="localhost",EMS_PORT="5001",MARKET_SIM_HOST="localhost",MARKET_SIM_PORT="5000",POV_PERCENTAGE="10",POV_MIN_SLICE="10",POV_MAX_SLICE="5000",OBSERVABILITY_HOST="localhost",OBSERVABILITY_PORT="5007",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vwap-algo]
command=deno run --allow-all /app/backend/src/algo/vwap-strategy.ts
environment=VWAP_ALGO_PORT="5006",EMS_HOST="localhost",EMS_PORT="5001",MARKET_SIM_HOST="localhost",MARKET_SIM_PORT="5000",VWAP_WINDOW_TICKS="20",OBSERVABILITY_HOST="localhost",OBSERVABILITY_PORT="5007",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:user-service]
command=deno run --allow-all /app/backend/src/user-service/user-service.ts
environment=USER_SERVICE_PORT="5008",USER_DB_PATH="/app/backend/data/users.db",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:journal]
command=deno run --allow-all /app/backend/src/journal/journal-server.ts
environment=JOURNAL_PORT="5009",JOURNAL_DB_PATH="/app/backend/data/journal.db",JOURNAL_RETENTION_DAYS="90",REDPANDA_BROKERS="localhost:9092"
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=deno run --allow-net --allow-read jsr:@std/http/file-server /app/frontend/dist --port 8080 --host 0.0.0.0
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
